# dev-util/git/files/git-daemon.initd
#!/sbin/runscript
# Copyright 1999-2008 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/dev-util/git/files/git-daemon.initd,v 1.3 2008/09/16 17:52:15 robbat2 Exp $

depend() {
	need net
	use logger
}

start() {
	ebegin "Starting git-daemon"
		start-stop-daemon --start --quiet --background \
		--exec /usr/bin/git -- daemon ${GITDAEMON_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping git-daemon"
		start-stop-daemon --stop --quiet --name git-daemon
	eend $?
}
