# app-misc/linux-logo/files/linux-logo.initscript
#!/sbin/runscript

start() {
	ebegin "Starting linux_logo"

	if [ ! -x /usr/bin/linux_logo ]
	then
		eerror "ERROR:  linux_logo not found !"
		return 1
	fi

	ebegin "  Creating /etc/issue"
	/usr/bin/linux_logo $LOGO $OPTIONS -F "${FORMAT}" > /etc/issue
	eend $? "  Failed to create /etc/issue"

	ebegin "  Creating /etc/issue.net"
	/usr/bin/linux_logo $LOGO $OPTIONS -F "${FORMATNET}" > /etc/issue.net
	eend $? "  Failed to create /etc/issue.net"
}
