# sys-apps/man/files/makewhatis.cron
#!/bin/sh
# this is part of the man package
# it updates the searchly database
# for manpages

/bin/nice /usr/sbin/makewhatis
